# systemd/battery-notification.service
[Unit]
Description=Battery Low Notification

[Service]
ExecStart=/bin/bash /home/misha/scripts/battery_check.sh

